# config/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/var/www/manual_providers/release/config/start_cron.sh
